# parsedmarc/crontab
*/10 * * * * parsedmarc -c /parsedmarc.ini --debug >> /var/log/cron.log 2>&1
